# voice_mode/templates/systemd/voicemode-broker.service
[Unit]
Description=VoiceMode hands-free conversation broker
After=default.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
ExecStart={START_SCRIPT}
Restart=always
RestartSec=5
Environment="PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin"
StandardOutput=journal
StandardError=journal
SyslogIdentifier=voicemode-broker

[Install]
WantedBy=default.target
